# systemd/getty@tty.service.d/skip-username.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -- arch' --noclear --skip-login - $TERM
