# applications/autostart/fusuma.desktop
[Desktop Entry]
Type=Application
Name=fusuma
Exec=fusuma -d
Comment=fusuma
