# applications/autostart/plank.desktop
[Desktop Entry]
Type=Application
Name=plank
Exec=plank 
Comment=plank
