# applications/autostart/setxkbmap.desktop
[Desktop Entry]
Type=Application
Name=setxkbmap
Exec=setxkbmap -option caps:escape
Comment=setxkbmap

