# applications/autostart/sxhkd-kde.desktop
[Desktop Entry]
Type=Application
Name=sxhkd-kde
Exec=/home/arch/Scripts/sxhkd-kde.sh
Comment=sxhkd
