# applications/autostart/warpd.desktop
[Desktop Entry]
Type=Application
Name=warpd
Exec=warpd
Comment=warpd
