# applications/autostart/xmodmap.desktop
[Desktop Entry]
Type=Application
Name=xmodmap
Exec=xmodmap ~/Dotfiles/X11/Xmodmap
Comment=xmodmap
