# applications/autostart/xremap.desktop
[Desktop Entry]
Type=Application
Name=xremap
Exec=xremap --watch /home/arch/.config/xremap/config.yml
Comment=xremap
