# kde-plasma/nordic/configs/autostart/sxhkd.desktop
[Desktop Entry]
Type=Application
Name=sxhkd
Exec=sxhkd 
Comment=sxhkd
