# kde-plasma/nordic/configs/autostart/xrdb.desktop
[Desktop Entry]
Type=Application
Name=X resources loader
Exec=xrdb -merge $HOME/Dotfiles/X11/Xresources 
Comment=xresources
